# deploy/crag-cleanup.timer.example
[Unit]
Description=Run cRAG results cleanup daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target

